# dployco-unicorn/templates/default/unicorn.conf.erb
# Dynamically generated by Chef on <%= node["fqdn"] %>
# Local modifications will be overwritten by Chef.

description "<%= @app_name %>"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
  <% if @socket.include? "/" -%>
    mkdir -p <%= File.dirname(@socket) %>
    chown -R <%= @owner %>:<%= @group %> <%= File.dirname(@socket) %>
  <% end -%>
end script

respawn
exec start-stop-daemon --start --chuid <%= @owner %>:<%= @group %> --chdir <%= @app_path %>/current --exec <%= @bundler ? "#{@bundle_command} exec" : '' %> <%= @smells_like_rack ? 'unicorn' : 'unicorn_rails' %> -- -E <%= @rails_env %> -c /etc/unicorn/<%= @app_name %>.rb